# root/etc/crontabs/abc
2,7,12,17,22,27,32,37,42,47,52,57 * * * * /app/duck.sh 2>&1
