# root/defaults/crontabs/abc
# min   hour    day     month   weekday command
*/5     *       *       *       *       sleep $((60 + $RANDOM % 230)); /app/duck.sh 2>&1
